# srcs/php_service/srcs_php/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php]
command=sh -c "php-fpm7 --nodaemonize && kill -s SIGTERM $(cat supervisord.pid)"

[program:nginx]
command=sh -c "/usr/sbin/nginx -g 'daemon off;' && kill -s SIGTERM $(cat supervisord.pid)"
