# srcs/telegraf_service/srcs_telegraf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:telegraf]
command=/bin/sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"
